# playbooks/roles/ff.iptables/templates/iptables.service.j2
[Unit]
Description=Set NAT Masquerade in iptables
After=network-online.target

[Service]
Type=oneshot

ExecStart=echo "configuring NAT if needed"
{% for interface in networkd_interfaces | default([]) %}
{% if interface.nat4 | default(False) %}
ExecStart=iptables -A POSTROUTING -o {{ interface.name }} -j MASQUERADE -t nat
{% endif %}
{% if interface.nat6 | default(False) %}
ExecStart=ip6tables -t nat -A POSTROUTING -o {{ interface.name }} -j MASQUERADE
{% endif %}
{% endfor %}

[Install]
WantedBy=multi-user.target